# templates/default/elite-golf-time.timer.erb
[Unit]
Description=Elite Golf Server Time Update Timer
Requires=elite-golf-time.service

[Timer]
OnActiveSec=1s
OnUnitActiveSec=2s

[Install]
WantedBy=timers.target